# docs/deploy/stackchan-gateway.service.d/heartbeat.conf
# Phase D1: heartbeat の opt-in（yorishiro fork 専用）
# 低頻度・無音の仕草（見回し・表情・うなずき）で「生きてる感」を出す。
# このファイルを外せば heartbeat は完全に無効（デフォルト挙動）に戻る。
# install: sudo cp <this file> /etc/systemd/system/stackchan-gateway.service.d/
#          sudo systemctl daemon-reload && sudo systemctl restart stackchan-gateway
[Service]
# 平均間隔（分）。±25% のジッター付き。検証時は 1 に下げると約1分間隔になる
Environment=STACKCHAN_HEARTBEAT_INTERVAL_MIN=30
# 任意（デフォルト値あり）:
# クワイエットアワー（ローカル時刻、跨ぎ可）。デフォルト 22:00-08:00。"off" で無効
#Environment=STACKCHAN_HEARTBEAT_QUIET=22:00-08:00
# 間隔ジッター 0..0.9（デフォルト 0.25）
#Environment=STACKCHAN_HEARTBEAT_JITTER=0.25
