# docs/deploy/stackchan-gateway.service.d/local-llm.conf
# Phase C3: 応答ルーティング層の opt-in（yorishiro fork 専用）
# 短い発話をローカル Ollama モデルに直行させる。このファイルを外せば従来挙動に戻る。
# install: sudo cp <this file> /etc/systemd/system/stackchan-gateway.service.d/
#          sudo systemctl daemon-reload && sudo systemctl restart stackchan-gateway
[Service]
Environment=STACKCHAN_LOCAL_LLM_MODEL=hf.co/LiquidAI/LFM2.5-1.2B-JP-202606-GGUF:Q4_K_M
# Phase D3: 家庭利用では会話間隔が 30 分を超えて毎回コールドロード(~3s)に
# なっていたため、常駐(0.7GB、GPU の他利用者なし)で warm 0.5s を維持する
Environment=STACKCHAN_LOCAL_LLM_KEEP_ALIVE=24h
# 任意（デフォルト値あり）:
#Environment=STACKCHAN_LOCAL_LLM_URL=http://127.0.0.1:11434
#Environment=STACKCHAN_LOCAL_LLM_TIMEOUT_S=10
